# infra/ec2-devbox/systemd/devbox-idle-stop.timer
[Unit]
Description=Check devbox idleness every 5 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=devbox-idle-stop.service

[Install]
WantedBy=timers.target
